# apps/stockpilot/systemd/stockpilot-production-daemon.service
[Unit]
Description=StockPilot Production Management Daemon
Documentation=https://stockpilot.ai/docs/production
After=network.target redis.service postgresql.service nginx.service
Wants=redis.service postgresql.service nginx.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=stockpilot
Group=stockpilot
WorkingDirectory=/opt/stockpilot/backend
Environment=PATH=/opt/stockpilot/backend/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/stockpilot/backend
Environment=LOG_LEVEL=INFO
EnvironmentFile=-/opt/stockpilot/.env.production
ExecStart=/opt/stockpilot/backend/venv/bin/python /opt/stockpilot/backend/production_daemon.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stockpilot-production

# 보안 강화 설정
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/stockpilot /tmp /opt/stockpilot/data /opt/stockpilot/backups /etc/nginx/sites-available /etc/nginx/sites-enabled /etc/ssl /etc/letsencrypt
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# 리소스 제한
MemoryMax=1G
CPUQuota=75%
TasksMax=200

# SSL 인증서 관리를 위한 추가 권한
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target